# systemd/hybrid-dns-monitor.service
[Unit]
Description=Hybrid DNS Server Monitoring Service
Documentation=https://github.com/user/hybrid-dns-server
After=network.target postgresql.service bind9.service
Wants=postgresql.service bind9.service

[Service]
Type=exec
User=dns-admin
Group=dns-admin
WorkingDirectory=/opt/hybrid-dns-server/monitoring

# Environment
EnvironmentFile=/opt/hybrid-dns-server/.env

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/var/log/bind
ReadWritePaths=/var/log/hybrid-dns

# Resource limits
LimitNOFILE=1024
TasksMax=100

# Start command
ExecStart=/opt/hybrid-dns-server/backend/venv/bin/python /opt/hybrid-dns-server/monitoring/simple_monitor.py
KillMode=mixed
Restart=on-failure
RestartSec=10
TimeoutStopSec=15

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hybrid-dns-monitoring

[Install]
WantedBy=multi-user.target